# greenhouse-app/deploy/greenhouse-monitor.service
[Unit]
Description=Greenhouse — monitor senzori (DHT11, ploaie, foc, gaz, sol, buzzer)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
User=greenhouse
# Acces la GPIO fara root (userul trebuie sa fie in grupul gpio):
#   sudo usermod -aG gpio greenhouse
SupplementaryGroups=gpio
WorkingDirectory=/home/greenhouse/Desktop/GreenHouse
ExecStart=/usr/bin/python3 greenhouse_monitor.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
